# mongo-bootstrapping-role/templates/services/mongod.service.d/mongod-set-pid-path.conf
# Informs systemd about the PID file path

[Service]
PIDFile={{ run_path }}/mongod.pid
